# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid
loglevel=info

; ── Bot Telegram ─────────────────────────────────────────────────────────────
[program:bot]
command=python bot.py
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── REST API (FastAPI + uvicorn) ─────────────────────────────────────────────
[program:api]
command=uvicorn api.main:app --host 0.0.0.0 --port 8000 --workers 1
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
